# servers/elk/cron.d/elk
# Run ELK curator
*/20 0 * * * root /bin/run-curator > /tmp/logs.curator
